# videocontrol.service
[Unit]
Description=VideoControl Server v3.0.0
After=network.target

[Service]
Type=simple
User=vcuser
Group=vcgroup
WorkingDirectory=/vid/videocontrol
ExecStart=/usr/bin/node /vid/videocontrol/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# ReadWritePaths настраиваются динамически в зависимости от contentRoot из настроек
# По умолчанию: /vid/videocontrol/data и /vid/videocontrol/config
# Если используется внешний диск, добавьте его путь сюда
ReadWritePaths=/vid/videocontrol/data /vid/videocontrol/config

# Environment
Environment=NODE_ENV=production
# Загружаем переменные из .env файла
EnvironmentFile=/vid/videocontrol/.env

[Install]
WantedBy=multi-user.target
