# videocontrol.service
[Unit]
Description=MMRC Server v3.2.0
After=network-online.target local-fs.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=vcuser
Group=vcgroup
WorkingDirectory=/var/lib/mmrc
ExecStartPre=/usr/bin/node /var/lib/mmrc/scripts/check-startup-compat.js
ExecStartPre=/usr/bin/node /var/lib/mmrc/src/database/migrate.js
ExecStart=/usr/bin/node /var/lib/mmrc/server.js
Restart=on-failure
RestartSec=3
TimeoutStartSec=120
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mmrc

# Security
NoNewPrivileges=true
PrivateTmp=true
# Позволяет менять contentRoot через админку без правки systemd unit.
# Защищает системные директории (/usr, /boot, /etc), но не ограничивает data-директорию.
ProtectSystem=full

# Resource limits
LimitNOFILE=65535

# Environment
Environment=NODE_ENV=production
Environment=HOME=/home/vcuser
Environment=ANDROID_USER_HOME=/home/vcuser/.android
Environment=ANDROID_SDK_HOME=/home/vcuser
# Загружаем переменные из .env файла
EnvironmentFile=-/var/lib/mmrc/.env

[Install]
WantedBy=multi-user.target
