# videocontrol.service
[Unit]
Description=VCServer
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
Group=vcgroup
WorkingDirectory=/vid/videocontrol
ExecStart=/usr/bin/node /vid/videocontrol/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/mnt/videocontrol-data /vid/videocontrol/data /vid/videocontrol/config

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
